# systemd/vcert.service
# =============================================================================
# systemd service (oneshot) to run the VCERT playbook
# Install at: /etc/systemd/system/vcert.service
# =============================================================================
[Unit]
Description=VCERT - certificate renewal (CyberArk Certificate Manager)
Documentation=https://github.com/Venafi/vcert
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# The token is read from the environment file (chmod 600, owned by root).
# Create /etc/vcert/vcert.env with:  VCERT_TOKEN=your_access_token
EnvironmentFile=/etc/vcert/vcert.env
ExecStart=/usr/local/bin/vcert run -f /etc/vcert/playbook.yaml
# Basic hardening
User=root
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
