# systemd/vcert.timer
# =============================================================================
# systemd timer - triggers vcert.service periodically
# Install at: /etc/systemd/system/vcert.timer
#
# Enable:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vcert.timer
#   systemctl list-timers vcert.timer
# =============================================================================
[Unit]
Description=Schedules certificate renewal via VCERT
Documentation=https://github.com/Venafi/vcert

[Timer]
# Runs twice a day. VCERT only renews when within the renewBefore window.
OnCalendar=*-*-* 06,18:00:00
# If the system was off at the scheduled time, run as soon as it boots.
Persistent=true
# Spread execution by up to 30min to avoid a simultaneous spike across large fleets.
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
